# packaging/deb/postinst
#!/usr/bin/env bash
set -euo pipefail

if ! id -u nuvion >/dev/null 2>&1; then
  useradd --system --no-create-home --shell /usr/sbin/nologin nuvion
fi

if getent group video >/dev/null 2>&1; then
  usermod -a -G video nuvion || true
fi

install -d -m 0770 -o root -g nuvion /etc/nuv-agent

if [ ! -f /etc/nuv-agent/agent.env ]; then
  cp /opt/nuv-agent/share/agent.env.example /etc/nuv-agent/agent.env
fi

chown root:nuvion /etc/nuv-agent/agent.env
chmod 0660 /etc/nuv-agent/agent.env

if [ -n "${SUDO_USER:-}" ] && [ "$SUDO_USER" != "root" ]; then
  usermod -a -G nuvion "$SUDO_USER" || true
  if getent group video >/dev/null 2>&1; then
    usermod -a -G video "$SUDO_USER" || true
  fi
  if command -v setfacl >/dev/null 2>&1; then
    setfacl -m "u:${SUDO_USER}:rwx" /etc/nuv-agent || true
    setfacl -m "u:${SUDO_USER}:rw" /etc/nuv-agent/agent.env || true
  fi
fi

if command -v python3 >/dev/null 2>&1; then
  python3 -m venv --clear --system-site-packages /opt/nuv-agent/venv
  /opt/nuv-agent/venv/bin/pip install --no-cache-dir --upgrade pip >/dev/null
  /opt/nuv-agent/venv/bin/pip install --no-cache-dir --upgrade /opt/nuv-agent/src >/dev/null
fi

systemctl daemon-reload || true
systemctl enable --now nuv-agent.service || true
